# terraform/app/files/docker-auto-update.timer
[Unit]
Description=Check for Docker image updates every 15 seconds

[Timer]
OnUnitActiveSec=15s
OnBootSec=15s

[Install]
WantedBy=timers.target
